# deploy/systemd/authd.service
[Unit]
Description=OpenTrusty Auth Plane (authd)
After=network.target postgresql.service

[Service]
Type=simple
User=opentrusty
Group=opentrusty
EnvironmentFile=/etc/opentrusty/authd.env
ExecStart=/usr/local/bin/authd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
